# systemd/n8n-monitoring.service
# N8N-R8 with Full Monitoring Stack Systemd Service File
# This service file starts N8N-R8 with complete monitoring (Prometheus, Grafana, etc.)
#
# INSTALLATION INSTRUCTIONS:
# ==========================
# 1. sudo cp /home/ram/projects/n8n-r8/systemd/n8n-monitoring.service /etc/systemd/system/
# 2. Edit paths and user in the service file
# 3. sudo systemctl daemon-reload
# 4. sudo systemctl enable n8n-monitoring.service
# 5. sudo systemctl start n8n-monitoring.service

[Unit]
Description=N8N-R8 with Full Monitoring Stack
Documentation=https://github.com/your-repo/n8n-r8
After=network-online.target docker.service
Wants=network-online.target
Requires=docker.service

[Service]
Type=forking
User=ram
Group=ram
WorkingDirectory=/home/ram/projects/n8n-r8

# Environment configuration
EnvironmentFile=/home/ram/projects/n8n-r8/.env
EnvironmentFile=-/home/ram/projects/n8n-r8/.env.local
Environment=COMPOSE_PROJECT_NAME=n8n-r8

# Service commands
ExecStartPre=-/usr/bin/docker compose -f docker-compose.yml -f docker-compose.monitoring.yml pull --quiet
ExecStartPre=-/usr/bin/docker compose -f docker-compose.yml -f docker-compose.monitoring.yml down --remove-orphans
ExecStartPre=/bin/bash -c 'mkdir -p data/{n8n,postgres,redis} monitoring/{data/{prometheus,grafana,alertmanager,loki,uptime-kuma},logs}'
ExecStartPre=/bin/bash -c 'chmod -R 755 data/ monitoring/'

ExecStart=/usr/bin/docker compose -f docker-compose.yml -f docker-compose.monitoring.yml up -d
ExecStop=/usr/bin/docker compose -f docker-compose.yml -f docker-compose.monitoring.yml down
ExecReload=/usr/bin/docker compose -f docker-compose.yml -f docker-compose.monitoring.yml restart

# Health check - wait longer for monitoring stack to be ready
ExecStartPost=/bin/sleep 60
ExecStartPost=/bin/bash -c 'curl -f http://localhost:5678/healthz && curl -f http://localhost:9090/-/healthy || exit 1'

# Restart policy
Restart=always
RestartSec=20
TimeoutStartSec=600  # Longer timeout for monitoring stack
TimeoutStopSec=120

# Resource limits for monitoring stack
LimitNOFILE=65536
MemoryLimit=6G

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=n8n-r8-monitoring

# Security
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
